# files/ssmserver.conf
[program:ssmserver]
command = /opt/Supermicro/SSM/SSMServer/start_server.sh
directory = /opt/Supermicro/SSM/SSMServer
autostart = true
autorestart = true
stopwaitsecs = 5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
